# .config/systemd/user/blueman-applet.service
[Unit]
Description=Blueman applet

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=desktop.target
