# .config/systemd/user/desktop.target
[Unit]
Description=Desktop environment
Wants=environment.target
Wants=window-manager.service
Wants=sound.target

[Install]
Alias=default.target
