# .config/systemd/user/display-set.target
[Unit]
Description=Target for services that are binded to specific monitor configuration
