# .config/systemd/user/dropboxd.service
[Unit]
Description=Dropbox

[Service]
ExecStart=/usr/bin/dropboxd

[Install]
WantedBy=desktop.target
