# .config/systemd/user/dropboxd.service
[Unit]
Description=Dropbox

[Service]
ExecStart=/usr/bin/dropboxd
Type=forking
SuccessExitStatus=0 SIGKILL

[Install]
WantedBy=desktop.target
