# .config/systemd/user/lid-close.target
[Unit]
Description=Lid close target
BindsTo=screen-lock.service
