# .config/systemd/user/sound.target
[Unit]
Description=Sound initialization target
