# .config/systemd/user/volumeicon.service
[Unit]
Description=Volume icon
Wants=sound.target
After=sound.target

[Service]
ExecStart=/usr/bin/volumeicon

[Install]
WantedBy=desktop.target
