# .config/systemd/user/xrr-events.service
[Unit]
Description=XrandR events handling daemon

[Service]
ExecStart=/usr/bin/xrr-events

[Install]
WantedBy=desktop.target
