# systemd/desktop.target
[Unit]
Description=Desktop environment
Wants=window-manager.service

[Install]
Alias=default.target
