# systemd/desktop.target
[Unit]
Description=Desktop environment
Wants=window-manager.service
Wants=sound.target

[Install]
Alias=default.target
