# systemd/desktop-norr.target
[Unit]
Description=Desktop environment parts without XRandR support
