# systemd/environment.target
[Unit]
Description=User environment variables target
