# systemd/lid-close.target
[Unit]
Description=Lid close target
Wants=lockscreen.target
