# systemd/nm-applet.service
[Unit]
Description=NetworkManager applet
Requires=NetworkManager.service

[Service]
ExecStart=/usr/bin/nm-applet

[Install]
WantedBy=desktop.target
