# systemd/volumeicon.service
[Unit]
Description=Volume icon

[Service]
ExecStart=/usr/bin/volumeicon

[Install]
WantedBy=desktop.target
