# systemd/xmenud.service
[Unit]
Description=xmenud applications menu

[Service]
ExecStart=/usr/bin/xmenud -t

[Install]
WantedBy=desktop.target
