# systemd/xxkb.service
[Unit]
Description=XXKB
After=layout.service

[Service]
ExecStart=/usr/bin/xxkb

[Install]
WantedBy=desktop.target
